# home/.config/bash/profile.d/env.sh
#
# Editors
#
export EDITOR='vim'
export VISUAL='vim'
export PAGER='less'

#
# Language
#
if test -z "$LANG"; then
    export LANG='ja_JP.UTF-8'
fi
